# packaging/debian/postinst
#!/bin/sh
set -e

case "$1" in
  configure)
    if ! getent passwd muros >/dev/null; then
      adduser --system --group --no-create-home --home /var/lib/muros muros
    fi

    # Application AND SSH administrator account (PAM authentication).
    # MurOS authenticates the web UI through PAM against a local Linux
    # account, and that same account is used for SSH: one set of
    # credentials covers both the UI and the shell. The default
    # administrator is the system 'root' account directly, so there is
    # no separate 'admin' user to create here, and root keeps its
    # existing password (MurOS never resets it). Every other Linux
    # account is refused entry to the web UI until root grants it access
    # from the Access > Users page (ui_access gate in the backend).

    install -d -m 750 -o muros -g muros /etc/muros
    install -d -m 750 -o muros -g muros /var/lib/muros
    install -d -m 750 -o muros -g muros /var/lib/muros/backups
    install -d -m 750 -o muros -g muros /var/cache/muros/upgrade
    # If install.sh restored a snapshot from /var/backups/muros into
    # /var/lib/muros, the muros UID/GID may not match the new package's
    # muros user. Force a recursive chown so the backend (running as
    # the muros user) can read its DB and write to /var/lib/muros/backups.
    chown -R muros:muros /var/lib/muros 2>/dev/null || true

    if [ ! -f /etc/muros/jwt.key ]; then
      umask 077
      openssl rand -hex 32 > /etc/muros/jwt.key
      chown muros:muros /etc/muros/jwt.key
      chmod 600 /etc/muros/jwt.key
    fi

    # Construction du venv directement sur la cible. C'est plus fiable
    # que de shipper un venv pre-rempli depuis le runner CI (shebang
    # casse, paths du builder, etc.).
    if [ ! -x /opt/muros/backend/venv/bin/uvicorn ]; then
      echo "  -> creation du venv Python (peut prendre une minute)"
      python3 -m venv /opt/muros/backend/venv
      /opt/muros/backend/venv/bin/pip install --quiet --upgrade pip wheel
      /opt/muros/backend/venv/bin/pip install --quiet -r /opt/muros/backend/requirements.txt
    fi

    # Capacites reseau pour python (netlink, raw sockets)
    PYBIN=$(readlink -f /opt/muros/backend/venv/bin/python3)
    /usr/sbin/setcap cap_net_admin,cap_net_raw=eip "$PYBIN" || true

    # Pas d'Alembic : le schema est cree au demarrage du backend via
    # SQLAlchemy Base.metadata.create_all(). Les nouvelles colonnes
    # ajoutees au fil des versions sont detectees par un script de
    # migration soft dans app/db.py (ALTER TABLE ADD COLUMN si manque).

    # Cert TLS : snakeoil Debian standard (livre par le paquet ssl-cert
    # liste en Depends). Sur images minimalistes le couple peut ne pas etre
    # genere automatiquement, on force ici.
    if [ ! -s /etc/ssl/certs/ssl-cert-snakeoil.pem ] \
       || [ ! -s /etc/ssl/private/ssl-cert-snakeoil.key ]; then
      make-ssl-cert generate-default-snakeoil --force-overwrite || true
    fi
    install -d -m 0755 /etc/nginx/ssl
    ln -sf /etc/ssl/certs/ssl-cert-snakeoil.pem /etc/nginx/ssl/muros.crt
    ln -sf /etc/ssl/private/ssl-cert-snakeoil.key /etc/nginx/ssl/muros.key

    if [ ! -L /etc/nginx/sites-enabled/muros ]; then
      ln -sf /etc/nginx/sites-available/muros /etc/nginx/sites-enabled/muros
      rm -f /etc/nginx/sites-enabled/default
    fi

    # Recharge sysctl pour activer ip_forward + hardening immediatement
    sysctl --system >/dev/null 2>&1 || true

    # Silence systemd-ssh-generator on bare-metal / VMs without VSOCK.
    # The generator (shipped by openssh-server) runs on every
    # `systemctl daemon-reload` and tries to query the local AF_VSOCK
    # CID to build a sshd@.service template for SSH-over-VSOCK. On any
    # machine without a virtio-vsock device the syscall fails with
    # EADDRNOTAVAIL and dumps "Failed to query local AF_VSOCK CID:
    # Cannot assign requested address" into the journal at every apply.
    # Since MurOS triggers daemon-reload on every configuration push,
    # the log gets unreadable fast. Masking the generator by symlinking
    # it to /dev/null is the upstream-recommended way (cf systemd man
    # page systemd.generator(7)). The real ssh.service unit is
    # unaffected.
    install -d -m 0755 /etc/systemd/system-generators
    if [ ! -L /etc/systemd/system-generators/systemd-ssh-generator ]; then
      ln -sf /dev/null /etc/systemd/system-generators/systemd-ssh-generator
    fi

    # MurOS = seul control plane reseau, comme FortiOS / OpenWrt / VyOS.
    # NetworkManager / ifupdown / netplan / resolvconf sont uninstalles
    # avant l'arrivee ici par le install.sh (et le Conflicts: dans
    # debian/control empeche dpkg de les coinstaller de toute facon).
    #
    # Reste systemd-networkd qui fait partie integrante du paquet
    # systemd (impossible a uninstall) : on le mask pour qu'il ne demarre
    # pas et ne ramene pas d'IPs depuis /etc/systemd/network. Idem pour
    # systemd-networkd-wait-online qui bloquerait le boot a attendre une
    # IP que personne ne va lui poser.
    #
    # systemd-resolved : meme histoire (livre par systemd), si actif il
    # remplace /etc/resolv.conf par un symlink vers son stub 127.0.0.53
    # et toute conf DNS MurOS est ignoree.
    #
    # dhcpcd : sur Debian 13 trixie, dhcpcd-base est le client DHCP par
    # defaut (remplace isc-dhcp-client). Son hook 20-resolv.conf reecrit
    # /etc/resolv.conf a chaque renew, ecrasant la conf MurOS. On stoppe
    # et on masque ; muros-boot prend le relais pour maintenir les IPs.
    # networking.service (ifupdown) is listed here too : on an ISO build,
    # debian-installer often installs `ifupdown` despite our Conflicts:
    # because preseed ordering varies. If left enabled, networking.service
    # runs at boot, spawns dhclient on the primary NIC and adds a parasite
    # default route with metric 1002 (= 1000 + ifindex) that fights with
    # muros-boot's route. Masking the unit is the symmetric counterpart
    # to install.sh purging the package on bare-metal installs.
    for svc in systemd-networkd.service systemd-networkd-wait-online.service \
               systemd-resolved.service dhcpcd.service dhclient.service \
               networking.service; do
      if systemctl list-unit-files "$svc" >/dev/null 2>&1; then
        systemctl stop "$svc" 2>/dev/null || true
        systemctl disable "$svc" 2>/dev/null || true
        systemctl mask "$svc" 2>/dev/null || true
      fi
    done

    # On an ISO build, debian-installer may have installed ifupdown as
    # priority=important even though debian/control declares Conflicts:.
    # Purge it now so it cannot come back via a future apt upgrade and
    # cannot drop config files in /etc/network/interfaces.d/ that a
    # manual `ifup` could replay.
    if dpkg-query -W -f='${Status}' ifupdown 2>/dev/null \
         | grep -q "install ok installed"; then
      DEBIAN_FRONTEND=noninteractive apt-get purge -y ifupdown \
        >/dev/null 2>&1 || true
    fi

    # systemd-resolved is a separate package on Debian 13 (trixie) and
    # later. Mask alone is not enough : an admin who runs
    # `systemctl unmask systemd-resolved` (typical mistake on a fresh
    # checkout) wakes it up and it grabs ::1:53, which then makes
    # Unbound fail to start with "Address already in use". Purge the
    # package so it can't come back unless explicitly reinstalled.
    if dpkg-query -W -f='${Status}' systemd-resolved 2>/dev/null \
         | grep -q "install ok installed"; then
      DEBIAN_FRONTEND=noninteractive apt-get purge -y systemd-resolved \
        >/dev/null 2>&1 || true
    fi
    # Tue les daemons DHCP encore en cours (hors systemd) qui pourraient
    # reecrire /etc/resolv.conf juste apres notre write.
    pkill -x dhcpcd 2>/dev/null || true
    pkill -x dhclient 2>/dev/null || true

    # Nettoyage des fichiers intermediaires legacy (rc <= 29). MurOS
    # n'utilise plus AUCUN fichier de config reseau distribution : tout
    # passe par muros-boot + commandes `ip` direct, DB SQLite = source
    # unique de verite.
    rm -f /etc/network/interfaces.d/muros /etc/muros/routes.conf 2>/dev/null || true

    # Sauvegarde puis neutralise /etc/network/interfaces pour eviter qu'un
    # admin qui essaierait un `ifup` manuel s'auto-flingue. Garde une
    # copie au cas ou il veuille relire ce qui y etait.
    if [ -f /etc/network/interfaces ] && [ ! -f /etc/network/interfaces.muros-bak ]; then
      cp /etc/network/interfaces /etc/network/interfaces.muros-bak
      cat > /etc/network/interfaces <<'EOF'
# Fichier neutralise par MurOS (sauvegarde dans interfaces.muros-bak).
# Les interfaces sont gerees depuis l'UI MurOS, pas ici.
# Garder seulement la loopback pour la coherence systeme.

auto lo
iface lo inet loopback
EOF
    fi

    # NTP : MurOS uses chrony as the time daemon (enabled by default).
    # chrony is more capable than systemd-timesyncd (real NTP server mode,
    # better step/slew control, chronyc introspection) and is the service
    # the UI Time page manages via /etc/chrony/conf.d/muros.conf. Mask
    # systemd-timesyncd so the two cannot fight over the system clock,
    # then enable chrony now and at boot.
    if systemctl list-unit-files systemd-timesyncd.service >/dev/null 2>&1; then
      systemctl stop systemd-timesyncd.service 2>/dev/null || true
      systemctl disable systemd-timesyncd.service 2>/dev/null || true
      systemctl mask systemd-timesyncd.service 2>/dev/null || true
    fi
    if systemctl list-unit-files chrony.service >/dev/null 2>&1; then
      systemctl enable --now chrony.service 2>/dev/null || true
    elif systemctl list-unit-files chronyd.service >/dev/null 2>&1; then
      systemctl enable --now chronyd.service 2>/dev/null || true
    fi

    # DNS : MurOS gere /etc/resolv.conf directement (Debian 13 base = pas
    # de systemd-resolved, pas de resolvconf). MurOS a deja maske
    # NetworkManager / systemd-networkd / ifupdown plus haut, donc
    # personne d'autre n'ecrit ce fichier. A l'install on pose
    # immediatement des resolveurs publics fiables (1.1.1.1, 4.4.4.4,
    # 8.8.8.8) sinon l'appliance peut se retrouver sans DNS si le DHCP
    # du LAN d'install est ensuite remplace (cas typique : install en
    # bridged puis bascule en NAT, ou DNS de la box plus joignable apres
    # cablage). L'admin peut changer depuis Systeme > DNS.
    # Backup de l'ancien resolv.conf au premier passage (idempotent via
    # le marker .muros-bak). Si c'est un symlink (cas systemd-resolved
    # stub), on saute le backup (rien d'utile a sauvegarder) et on
    # ecrase quand meme le fichier. On ecrase a chaque install meme si
    # un fichier MurOS existe deja : les DHCP clients tiers peuvent
    # avoir reecrit le fichier entre-temps, l'admin reconfigurera depuis
    # l'UI si besoin.
    if [ -f /etc/resolv.conf ] && [ ! -L /etc/resolv.conf ] \
         && [ ! -f /etc/resolv.conf.muros-bak ]; then
      cp /etc/resolv.conf /etc/resolv.conf.muros-bak 2>/dev/null || true
    fi
    rm -f /etc/resolv.conf
    cat > /etc/resolv.conf <<'EOF'
# Generated by MurOS at install - edit from the UI System > DNS.
# Default public resolvers (reachable without specific LAN config).
nameserver 1.1.1.1
nameserver 4.4.4.4
nameserver 8.8.8.8
EOF
    chmod 0644 /etc/resolv.conf
    # Empeche un futur dhcpcd/dhclient de toucher resolv.conf au cas ou
    # il aurait survecu au mask plus haut (binaire lance manuellement,
    # etc.). chattr +i n'est pas portable, on prefere neutraliser les
    # hooks dhcpcd/dhclient si presents.
    for hook in /lib/dhcpcd/dhcpcd-hooks/20-resolv.conf \
                /etc/dhcp/dhclient-enter-hooks.d/resolved \
                /etc/dhcp/dhclient-exit-hooks.d/debian; do
      if [ -f "$hook" ] && [ ! -f "$hook.muros-disabled" ]; then
        mv "$hook" "$hook.muros-disabled" 2>/dev/null || true
      fi
    done

    # Desactive par defaut les services qui exigent une configuration
    # explicite avant d'avoir un sens (HA, VPN, IPsec). Tant que l'admin
    # ne les a pas configures depuis l'UI ils n'ont rien a faire et leur
    # etat parasite le dashboard. L'UI les re-active a l'apply.
    # The block below resets MurOS-managed services to "off until the
    # admin enables it from the UI". This is the correct default for a
    # FRESH install but absolutely wrong on a package upgrade: it would
    # stop the DNS resolver, the DHCP server, WireGuard, IPsec,
    # keepalived and friends every time the operator clicks "Install
    # update" in the UI, exactly the services they need to stay up.
    # Gate the reset on $2 being empty (no previous version installed).
    if [ -z "$2" ]; then
      # The administrator account is the system 'root' account and MurOS
      # keeps its existing password untouched: the operator logs into the
      # web UI with the same root password they already use for the shell
      # / console. We deliberately do NOT reset it nor force a change.

      # SSH closed by default. openssh-server is a hard dependency so the
      # binaries are present, but the daemon does not run until the admin
      # explicitly enables it from the UI (SSH access page) once keys and
      # zone restrictions are in place. On Debian 13 (trixie) sshd is
      # socket-activated, so we disable both ssh.service and ssh.socket.
      for sshunit in ssh.service ssh.socket sshd.service; do
        if systemctl list-unit-files "$sshunit" >/dev/null 2>&1; then
          systemctl disable --now "$sshunit" 2>/dev/null || true
          systemctl reset-failed "$sshunit" 2>/dev/null || true
        fi
      done

      # Services that need no prior configuration are enabled by default.
      # fail2ban protects the web UI / SSH out of the box; it has sane
      # built-in jails and our /etc/fail2ban/jail.d/muros.local drop-in.
      # (nginx, muros-backend, muros-boot and chrony are enabled
      # above / by dh_installsystemd; snmpd is brought up by the
      # backend at boot from its DB config.)
      if systemctl list-unit-files fail2ban.service >/dev/null 2>&1; then
        systemctl enable --now fail2ban.service 2>/dev/null || true
      fi

      # Helper: disable a feature daemon AND clear any lingering "failed"
      # status. A plain `disable --now` stops and un-enables the unit but
      # leaves it marked "failed" in systemd if it had crashed earlier
      # (e.g. a dependency daemon shipped by us, like dnsmasq or unbound,
      # was already installed and auto-started by its own package
      # postinst, collided on port 53, and crashed before our policy-rc.d
      # could hold it back). Without `reset-failed` the operator sees a
      # red "failed" unit right after install even though the feature is
      # simply not enabled yet. reset-failed makes the post-install state
      # clean: every optional daemon is "inactive (dead)", never "failed".
      hold_feature_daemon() {
        _svc="$1"
        if systemctl list-unit-files "$_svc" >/dev/null 2>&1; then
          systemctl disable --now "$_svc" 2>/dev/null || true
          systemctl reset-failed "$_svc" 2>/dev/null || true
        fi
      }

      for svc in keepalived.service conntrackd.service strongswan.service \
                 strongswan-starter.service; do
        hold_feature_daemon "$svc"
      done
      hold_feature_daemon wg-quick@wg0.service

      # Le watcher d'alertes (boucle SMTP/webhook) est inutile tant qu'aucun
      # canal de notif n'est configure. L'UI le start a l'apply notifications.
      hold_feature_daemon muros-watcher.service

      # V1.2 network services : Kea (DHCP) and Unbound (recursive DNS).
      # These run without any prior configuration, so MurOS enables them
      # by default (alongside NTP/chrony and SNMP). Kea is DHCP-only and
      # never binds port 53, so it coexists cleanly with Unbound: no
      # collision, both can stay up from install.
      #
      # Kea : the stock Debian config can reference a sample subnet that
      # fails to bind on a fresh box. Drop a minimal, always-valid MurOS
      # config first (no interface, no subnet = idle daemon), so the
      # service starts clean and hands out nothing until the admin
      # configures a pool from the UI (which rewrites this file).
      install -d -m 0755 /etc/kea
      if [ ! -f /etc/kea/kea-dhcp4.conf ] \
         || ! grep -q "managed by MurOS" /etc/kea/kea-dhcp4.conf; then
        cat > /etc/kea/kea-dhcp4.conf <<'KEAEOF'
// /etc/kea/kea-dhcp4.conf -- managed by MurOS, do not edit.
// Idle default config: no interface, no subnet. The MurOS backend
// rewrites this file when a DHCP pool is configured from the UI.
{
  "Dhcp4": {
    "interfaces-config": { "interfaces": [] },
    "control-socket": {
      "socket-type": "unix",
      "socket-name": "/run/kea/kea4-ctrl-socket"
    },
    "lease-database": {
      "type": "memfile",
      "persist": true,
      "name": "/var/lib/kea/kea-leases4.csv"
    },
    "valid-lifetime": 43200,
    "subnet4": [],
    "loggers": [{
      "name": "kea-dhcp4",
      "severity": "INFO",
      "output_options": [{ "output": "syslog" }]
    }]
  }
}
KEAEOF
      fi
      install -d -m 0750 /var/lib/kea
      if systemctl list-unit-files kea-dhcp4-server.service >/dev/null 2>&1; then
        systemctl unmask kea-dhcp4-server.service 2>/dev/null || true
        systemctl reset-failed kea-dhcp4-server.service 2>/dev/null || true
        systemctl enable --now kea-dhcp4-server.service 2>/dev/null || true
      fi
      if systemctl list-unit-files unbound.service >/dev/null 2>&1; then
        systemctl unmask unbound.service 2>/dev/null || true
        systemctl reset-failed unbound.service 2>/dev/null || true
        systemctl enable --now unbound.service 2>/dev/null || true
      fi

      # Multi-WAN monitor : inactif par defaut. Tant qu'aucun WAN gateway
      # n'est declare dans l'UI, il n'a rien a probe. L'admin l'enable
      # depuis la page WAN ; ici on s'assure juste qu'il n'est pas en
      # train de tourner avec une DB vide (waste CPU pour rien).
      hold_feature_daemon muros-wan-monitor.service
    fi

    # unattended-upgrades : applique la planification via le helper
    # /usr/lib/muros/apply-unattended.sh. Source de verite :
    # /etc/muros/unattended.json (gere par l'UI). Si absent, le helper
    # genere une planif par defaut basee sur hash(hostname) pour ecarter
    # les peers HA d'au moins 1 jour. Le paquet muros est blacklist (cf
    # 52muros-unattended) : ses upgrades passent toujours par l'UI.
    /usr/lib/muros/apply-unattended.sh || true

    # snmpd / fail2ban / nginx / muros-backend / muros-boot : activations
    # standard, deja gerees par dh_installsystemd. Aucune action a faire.

    systemctl daemon-reload
    # muros-watcher reste disabled tant qu'aucun canal SMTP/webhook n'est
    # configure depuis l'UI. Il sera enabled+start a l'apply notifications.
    systemctl enable muros-boot.service muros-backend.service
    # On lance muros-boot des l'install pour adopter la conf reseau
    # actuelle (interfaces + IPs + default route) dans la DB MurOS. Sans
    # ca, l'admin verrait la page Routes vide et au prochain reboot le
    # kernel rejouerait une DB vide (perte d'IP DHCP).
    # L'adoption est idempotente : si elle a deja eu lieu (marker
    # /var/lib/muros/.adopted), on skip.
    systemctl start muros-boot.service || true
    systemctl restart muros-backend.service || true
    systemctl reload nginx || systemctl restart nginx || true
    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;

  *)
    echo "postinst called with unknown argument $1" >&2
    exit 1
    ;;
esac

exit 0
