# debian/postinst
#!/bin/sh
set -e
case "$1" in
  configure)
    # Take over the php-fpm www pool with the MurOS-sized pool.
    if [ ! -e /etc/php/8.4/fpm/pool.d/www.conf.distrib ]; then
      dpkg-divert --package muros --add --rename \
        --divert /etc/php/8.4/fpm/pool.d/www.conf.distrib \
        /etc/php/8.4/fpm/pool.d/www.conf || true
    fi

    # Self-signed certificate for the web UI (replace later from the GUI).
    if [ ! -s /usr/local/etc/muros/server.crt ] || [ ! -s /usr/local/etc/muros/server.key ]; then
      mkdir -p /usr/local/etc/muros
      openssl req -x509 -newkey rsa:2048 -nodes -days 3650 \
        -keyout /usr/local/etc/muros/server.key \
        -out /usr/local/etc/muros/server.crt \
        -subj '/O=MurOS/CN=MurOS' \
        -addext 'basicConstraints=CA:FALSE' \
        -addext 'keyUsage=digitalSignature,keyEncipherment' \
        -addext 'extendedKeyUsage=serverAuth' \
        -addext 'subjectAltName=DNS:MurOS' >/dev/null 2>&1 || true
      chown www-data:www-data /usr/local/etc/muros/server.crt /usr/local/etc/muros/server.key 2>/dev/null || true
      chmod 600 /usr/local/etc/muros/server.key 2>/dev/null || true
    fi

    # Seed a default configuration so the control plane and UI can start.
    if [ ! -e /conf/config.xml ] && [ -e /usr/local/etc/config.xml.sample ]; then
      mkdir -p /conf
      cp -p /usr/local/etc/config.xml.sample /conf/config.xml
      chown www-data:www-data /conf/config.xml 2>/dev/null || true
      chmod 640 /conf/config.xml 2>/dev/null || true
    fi
    if [ ! -e /usr/local/etc/config.xml ]; then
      ln -sf /conf/config.xml /usr/local/etc/config.xml
    fi

    # Many ported OPNsense scripts carry a "#!/usr/local/bin/php" shebang (the
    # FreeBSD interpreter path) and are executed directly, e.g. the console menu
    # runs /usr/local/opnsense/scripts/shell/setports.php as a program. Debian
    # ships PHP at /usr/bin/php, so provide the expected path as a symlink.
    if [ ! -e /usr/local/bin/php ]; then
      mkdir -p /usr/local/bin
      ln -sf /usr/bin/php /usr/local/bin/php
    fi

    # MurOS runtime directories (/var/etc, PHP temp/session/cache). These are
    # declared in /usr/lib/tmpfiles.d/muros.conf and recreated at every boot;
    # create them now so the first configuration run after install succeeds
    # without waiting for a reboot.
    systemd-tmpfiles --create /usr/lib/tmpfiles.d/muros.conf >/dev/null 2>&1 || true
    # Fallback in case the tmpfiles file is unavailable for any reason.
    for d in /var/etc /var/lib/php/tmp /var/lib/php/sessions /var/lib/php/opcache; do
      mkdir -p "$d"
    done
    for d in /var/lib/php/tmp /var/lib/php/sessions; do
      chown www-data:www-data "$d" 2>/dev/null || true
      chmod 700 "$d" 2>/dev/null || true
    done
    # The CLI on-disk OPcache is written by root (configd) and read by www-data,
    # so it is group-owned by www-data with group write.
    chown root:www-data /var/lib/php/opcache 2>/dev/null || true
    chmod 775 /var/lib/php/opcache 2>/dev/null || true

    # Persistent, size-bounded systemd journal. The firewall log viewer reads
    # the kernel netfilter log lines from the journal, so it must survive reboots
    # and keep enough history while staying bounded on an appliance disk.
    if [ ! -e /etc/systemd/journald.conf.d/muros.conf ]; then
      mkdir -p /etc/systemd/journald.conf.d
      cat > /etc/systemd/journald.conf.d/muros.conf <<'JRN'
[Journal]
Storage=persistent
SystemMaxUse=512M
RuntimeMaxUse=64M
JRN
    fi
    mkdir -p /var/log/journal
    systemd-tmpfiles --create --prefix /var/log/journal >/dev/null 2>&1 || true
    systemctl restart systemd-journald 2>/dev/null || true

    # Kea DHCP is driven by MurOS: configd renders /etc/kea/*.conf from the
    # stored configuration and the kea_service.sh wrapper starts only the
    # daemons enabled in the UI. Keep the stock units from auto-starting at
    # boot with their packaged default configuration; MurOS owns their state.
    for u in kea-dhcp4-server kea-dhcp6-server kea-ctrl-agent kea-dhcp-ddns-server; do
      systemctl disable --now "$u" 2>/dev/null || true
    done

    # Unbound is driven by MurOS, not by the packaged systemd unit: configd
    # renders /var/unbound/unbound.conf and the unbound start.sh script launches
    # the daemon from there. Stop the stock unit so it does not hold port 53 with
    # its default /etc/unbound configuration, and keep it from starting at boot.
    systemctl disable --now unbound 2>/dev/null || true
    systemctl disable --now unbound-resolvconf 2>/dev/null || true
    # The packaged AppArmor profile confines unbound to /etc/unbound and
    # /var/lib/unbound, which would deny the MurOS working directory
    # /var/unbound. Disable the profile so unbound can run from there.
    if [ -e /etc/apparmor.d/usr.sbin.unbound ]; then
      mkdir -p /etc/apparmor.d/disable
      ln -sf /etc/apparmor.d/usr.sbin.unbound /etc/apparmor.d/disable/usr.sbin.unbound
      apparmor_parser -R /etc/apparmor.d/usr.sbin.unbound 2>/dev/null || true
    fi

    # strongSwan (charon-systemd) is driven by MurOS as well: configd renders
    # /etc/swanctl/swanctl.conf and the IPsec apply starts the daemon and loads
    # the connections itself. The packaged unit runs "swanctl --load-all" as an
    # ExecStartPost, which makes the whole daemon fail to start when a single
    # generated connection is invalid. Clear that hook with a drop-in so charon
    # always comes up and configuration errors only fail the explicit reload,
    # and keep the unit from auto-starting with an empty configuration at boot.
    install -d /etc/systemd/system/strongswan.service.d
    cat > /etc/systemd/system/strongswan.service.d/muros.conf <<'STRONGSWAN'
[Service]
ExecStartPost=
STRONGSWAN
    systemctl disable strongswan 2>/dev/null || true

    # Purge the compiled Volt template cache so the UI/template changes shipped
    # in this package are recompiled from the new sources on the next request,
    # instead of being served from stale compiled copies left by a prior version.
    rm -f /var/lib/php/cache/*.php 2>/dev/null || true
    # Purge the CLI on-disk OPcache for the same reason: drop any opcode files
    # compiled from the previous version so the new sources are picked up.
    rm -rf /var/lib/php/opcache/* 2>/dev/null || true

    # Pick up the lighttpd drop-in and (re)start the front end.
    systemctl daemon-reload 2>/dev/null || true
    systemctl reload-or-restart php8.4-fpm 2>/dev/null || true
    systemctl restart lighttpd 2>/dev/null || true

    # The configd control plane backs every dashboard widget and settings-apply
    # action over its local socket. Make sure it is enabled and running after
    # install/upgrade: if it is down the UI hangs on the connect timeout and all
    # widgets fail to load. It only listens on a unix socket, so this is safe.
    systemctl enable muros-configd 2>/dev/null || true
    systemctl restart muros-configd 2>/dev/null || true
    ;;
esac
#DEBHELPER#
exit 0
