# src/etc/systemd/system/muros-interfaces.service
# MurOS interface configuration (iproute2).
# Applies the addressing defined in /conf/config.xml to the assigned data
# interfaces. The management interface stays under systemd-networkd.
[Unit]
Description=MurOS interface configuration (iproute2)
Documentation=https://github.com/ecritel/muros
After=systemd-networkd.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
# Render the addressing of the assigned data interfaces into systemd-networkd
# units (DHCP for WAN, static for LAN), then run the iproute2 pass for any
# interface that networkd does not own. The management link is never touched.
ExecStart=/usr/bin/php /usr/local/opnsense/scripts/interfaces/render_networkd.php
ExecStart=/usr/bin/php /usr/local/opnsense/scripts/interfaces/apply_linux.php
ExecReload=/usr/bin/php /usr/local/opnsense/scripts/interfaces/render_networkd.php
ExecReload=/usr/bin/php /usr/local/opnsense/scripts/interfaces/apply_linux.php

[Install]
WantedBy=multi-user.target
