# src/etc/systemd/system/muros-configd.service
[Unit]
Description=MurOS configuration daemon (web UI control plane bridge)
Documentation=https://github.com/ecritel/muros
After=network.target

[Service]
# configd ships a FreeBSD-style double-fork daemonizer; on Debian we let systemd
# supervise it directly by running the built-in foreground "console" mode.
Type=simple
WorkingDirectory=/usr/local/opnsense/service
ExecStart=/usr/bin/python3 /usr/local/opnsense/service/configd.py console
ExecStartPre=/bin/rm -f /var/run/configd.socket
Restart=always
RestartSec=2
KillMode=control-group

[Install]
WantedBy=multi-user.target
